# dwmstatus.service
[Unit]
Description=Dwmstatus to file

[Service]
ExecStart=/usr/bin/dwmstatus-file

[Install]
WantedBy=multi-user.target